# overlays/base/etc/init.d/pinet
#!/sbin/openrc-run

# this will setup netoworking, like raspbian does
# using ssh & wpa_supplicant.conf file
# if those files are present, it will copy wpa_supplicant and start the services

name="pinet"

start() {
  /etc/pinet /media/mmcblk0p1
}