# overlays/branding/etc/profile.d/hello.sh
#!/bin/sh

echo "Howdy Hacker!"
cat /etc/radical_edward.ans
